# disabled_systemd/autostart-healthcheck.timer
[Unit]
Description=Run autostart health check every 24 hours

[Timer]
OnBootSec=24h
OnUnitActiveSec=24h
Persistent=true
Unit=autostart-healthcheck.service

[Install]
WantedBy=timers.target
