# systemd/autostart-periodic-status.service
[Unit]
Description=Periodic full status report for autostart stack
After=network-online.target autostart-vnc-keepalive.service autostart-test-bot.service autostart-tuna-backup.service
Wants=network-online.target

[Service]
Type=oneshot
User=user
Group=user
WorkingDirectory=/home/user/autostart
Environment=HOME=/home/user
ExecStart=/home/user/autostart/periodic_status.sh

[Install]
WantedBy=multi-user.target
