# systemd/autostart-tuna-backup.service.template
[Unit]
Description=Autostart tuna tunnel backup service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{AUTOSTART_USER}}
Group={{AUTOSTART_GROUP}}
WorkingDirectory={{AUTOSTART_DIR}}
Environment=HOME={{AUTOSTART_HOME}}
ExecStart={{AUTOSTART_DIR}}/tuna_headless.sh
Restart=always
RestartSec=15

[Install]
WantedBy=multi-user.target
